# manny-threatintel.service
[Unit]
Description=Manny Threat Intel Server
After=network.target

[Service]
Type=simple
User=REPLACE_WITH_USERNAME
WorkingDirectory=REPLACE_WITH_INSTALL_PATH
ExecStart=/usr/bin/python3 REPLACE_WITH_INSTALL_PATH/server.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=REPLACE_WITH_INSTALL_PATH

[Install]
WantedBy=multi-user.target
